# modules/etherpad/templates/upstart_etherpad.conf.erb
description "etherpad"

start on started networking
stop on runlevel [!2345]

env EPHOME=<%= @etherpad_dir %>
env EPLOGS=/var/log/etherpad
env EPUSER=<%= @etherpad_user %>

respawn

pre-start script
    cd $EPHOME
    mkdir -p $EPLOGS                           ||true
    chown $EPUSER:admin $EPLOGS                ||true
    chmod 0755 $EPLOGS                         ||true
    mkdir -p $EPHOME/var                       ||true
    chown -R $EPUSER:admin $EPHOME/var         ||true
    $EPHOME/bin/installDeps.sh >> $EPLOGS/error.log || { stop; exit 1; }
end script

script
  cd $EPHOME/
  exec su -s /bin/sh -c 'HOME=/tmp exec "$0" "$@"' $EPUSER -- /usr/local/<%= @global_nodejs_version %>/bin/node node_modules/ep_etherpad-lite/node/server.js \
                        >> $EPLOGS/access.log \
                        2>> $EPLOGS/error.log
end script
